# crypto-boiler.service
[Unit]
Description=Crypto Boiler Temperature Controller
After=network.target
Wants=network.target

[Service]
Type=simple
User=user
Group=user
WorkingDirectory=/home/user/CONTROLLER
ExecStart=/usr/bin/python3 /home/user/CONTROLLER/start_system.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Переменные окружения (при необходимости)
Environment=PYTHONPATH=/home/user/CONTROLLER
Environment=GPIO_PIN=17
Environment=ASIC_IP=192.168.0.127
Environment=MAX_TEMP=52.0
Environment=MIN_TEMP=51.9

[Install]
WantedBy=multi-user.target 